# docker/configs/apache-supervisor.conf
[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/apache.log
