# docker/configs/vite-watch.conf
[program:vite-watch]
command=npm run watch
directory=/var/www/html
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/vite-watch.log
user=www-data
